# deamon/tqr_service_ltc.conf
# Info
description "LTC Worker"
author      "Lou Carther"

# Events
start on startup
stop on shutdown

# Automatically respawn
respawn
respawn limit 20 5

# Run the script!
# Note, in this example, if your PHP script returns
# the string "ERROR", the daemon will stop itself.
script
    [ $(exec /usr/bin/php -f /var/www/vhosts/trenqr.com/dm1606/product/entities/mod.wamp/LtcRouterStart.php) = 'ERROR' ] && ( stop; exit 1; )
end script
